# supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command=/usr/sbin/sshd
autorestart=true
[program:smbd]
command=ionice -c 3 smbd -FS
autorestart=true
